# autostart/mystartup.desktop
[Desktop Entry]
Name=Startup Script
Comment=eiei
Type=Application
Exec=env mystartup.sh
